# daemon_nginx/xelis_daemon.service
[Unit]
Description=Xelis

[Service]
Type=forking
ExecStart=/root/xelis/target/release/xelis_daemon

[Install]
WantedBy=multi-user.target